# .config/yadm/etc/cron.daily/eclean
#!/usr/bin/env bash
sudo eclean-kernel --num 1
sudo eclean-dist
echo "Daily eclean ran!" >> /tmp/dailyEcleanFlag
