# .config/yadm/etc/cron.daily/flatpak-update
#!/usr/bin/env bash
sudo flatpak update -y
